# fleet/daily_worker.timer
[Unit]
Description=Start the daily worker-container

[Timer]
OnCalendar=*:0/10
#OnCalendar=02:00
Persistent=true
